# samba-server/supervisord.ini
[supervisord]
nodaemon=true
loglevel=info

# set some defaults and start samba in foreground (-F), logging to stdout (-S), and using our config (-s path)

[program:smbd]
command=smbd --no-process-group -F -S -s /etc/samba/smb.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nmbd]
command=nmbd --no-process-group -F -S -s /etc/samba/smb.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
