# rutorrent/rootfs/etc/service/rtorrent/run
#!/bin/sh

exec 2>&1
export TERM=xterm
exec su - torrent -c "/usr/bin/rtorrent"
